# program/sitiopc-agent.service
[Unit]
Description=SitioPC Agent
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/sitiopc-agent
ExecStart=/usr/bin/python3 /home/pi/sitiopc-agent/main.py
Restart=always
RestartSec=5
Environment="SITIOPC_SERVER=http://localhost:8081"
Environment="SITIOPC_POLL_INTERVAL=10"

[Install]
WantedBy=multi-user.target
